# contrib/init/navcoind.conf
description "HTS Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env HTSD_BIN="/usr/bin/HTSd"
env HTSD_USER="HTS"
env HTSD_GROUP="HTS"
env HTSD_PIDDIR="/var/run/HTSd"
# upstart can't handle variables constructed with other variables
env HTSD_PIDFILE="/var/run/HTSd/HTSd.pid"
env HTSD_CONFIGFILE="/etc/HTS/HTS.conf"
env HTSD_DATADIR="/var/lib/HTSd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # HTSd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$HTSD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run HTSd."
        echo "The setting must appear in $HTSD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"HTS Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$HTSD_PIDDIR"
    chmod 0755 "$HTSD_PIDDIR"
    chown $HTSD_USER:$HTSD_GROUP "$HTSD_PIDDIR"
    chown $HTSD_USER:$HTSD_GROUP "$HTSD_CONFIGFILE"
    chmod 0660 "$HTSD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$HTSD_PIDFILE" \
    --chuid $HTSD_USER:$HTSD_GROUP \
    --exec "$HTSD_BIN" \
    -- \
    -pid="$HTSD_PIDFILE" \
    -conf="$HTSD_CONFIGFILE" \
    -datadir="$HTSD_DATADIR" \
    -disablewallet \
    -daemon

